# modules/ei_bps/templates/wso2ei-business-process.service.erb
[Unit]
Description=WSO2 Enterprise Integrator - Business Process
After=network.target

[Service]
ExecStart=<%=@install_path%>/wso2/<%=@profile%>/bin/wso2server.sh start
ExecStop=<%=@install_path%>/wso2/<%=@profile%>/bin/wso2server.sh stop
ExecRestart=<%=@install_path%>/wso2/<%=@profile%>/bin/wso2server.sh restart
PIDFile=<%=@install_path%>/wso2/<%=@profile%>/wso2carbon.pid
User=<%=@user%>
Group=<%=@user_group%>
Type=forking
Restart=on-failure
RestartSec=5
StartLimitInterval=120s
StartLimitBurst=3

[Install]
WantedBy=multi-user.target
